# debian/wb-mqtt-w1.postinst
#!/bin/bash

deb-systemd-helper enable wb-mqtt-w1.service >/dev/null || true

#DEBHELPER#

